# backend/supervisor/conf.d/main.conf
[program:main]
directory=/app
command=python3 /app/main.py
autostart=true
autorestart=true
redirect_stderr=true
