# page/supervisor/conf.d/insert.conf
[program:insert]
directory=/app
command=python3 /app/insert.py
autostart=true
autorestart=true
redirect_stderr=true
